# deploy/catalogus.service
# systemd unit — install on VPS:
#   sudo cp deploy/catalogus.service /etc/systemd/system/catalogus.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable catalogus
#   sudo systemctl start catalogus
#
# Edit User= and WorkingDirectory= to match your server path.

[Unit]
Description=Catalogus Next.js (superclones.cloud)
After=network.target mariadb.service mysql.service
Wants=network-online.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/var/www/superclones.cloud
Environment=NODE_ENV=production
EnvironmentFile=/var/www/superclones.cloud/.env
# Port 3001 — port 3000 is used by inkoop-autos (PM2) on this VPS
ExecStart=/usr/bin/npm run start -- -p 3001
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=catalogus

[Install]
WantedBy=multi-user.target
